# configs/crontab
*/5 * * * * /usr/sbin/logrotate /etc/logrotate.conf --state /logs/logrotate-state